# systemd/.config/systemd/user/chime.service
[Unit]
Description=Quarterly Hour Chime Service

[Service]
Type=oneshot
ExecStart=/home/phnaharris/bin/chime

[Install]
WantedBy=default.target
